# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "pid /tmp/nginx.pid; daemon off;"

[program:jenkins]
user=jenkins
command=/usr/local/bin/jenkins.sh
